# memory/scripts/memory-guardian.service
[Unit]
Description=Memory Guardian - 记忆系统守护进程
After=network.target

[Service]
Type=simple
User=li
WorkingDirectory=/home/li/.openclaw/workspace/memory/scripts
ExecStart=/usr/bin/node /home/li/.openclaw/workspace/memory/scripts/memory-guardian.js
Restart=always
RestartSec=10

# 环境
Environment=NODE_ENV=production

# 日志
StandardOutput=append:/home/li/.openclaw/workspace/memory/logs/guardian.log
StandardError=append:/home/li/.openclaw/workspace/memory/logs/guardian-error.log

[Install]
WantedBy=multi-user.target
